# LedManager.service
[Unit]
Description=Manage leds connected with DMX protocol

[Service]
ExecStart=/usr/bin/python /home/ignapk/swiatelka/LedManager.py
